# proxy/yt-proxy.service
[Unit]
Description=YouTube/Media Proxy for RK322x
After=network.target

[Service]
User=USER_PLACEHOLDER
ExecStart=PYTHON_PLACEHOLDER PROXY_PATH_PLACEHOLDER
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
